# test-alpine/java.supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:mysql]
command=/usr/bin/mysqld_safe
user=mysql
priority=0

[program:redis]
command=/usr/bin/redis-server --protected-mode no